# etc/supervisord.d/ThesisPlanet.ini
# This file is autogenerated by Puppet. Manual changes will be overwritten!
[program:TP-DEP-QUEUE]
command=php /usr/share/TP-DEP/src/library/Gearman/Worker.php
numprocs=1
numprocs_start=0
priority=999
autostart=true
autorestart=true
startsecs=1
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
directory=/usr/share/TP-DEP/src/library/Gearman
user=web
redirect_stderr=false
stdout_logfile=/var/log/supervisor/TP-DEP-QUEUE/TP-DEP-QUEUE.out
stdout_logfile_maxbytes=250MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/TP-DEP-QUEUE/TP-DEP-QUEUE.err
stderr_logfile_maxbytes=250MB
stderr_logfile_backups=10
